# etc/supervisor/mwop.net.conf
[program:mwop_net]
directory=%release_path%
command=php7.4 vendor/bin/mezzio-swoole start
redirect_stderr=true
autostart=true
autorestart=true
user=www-data
killasgroup=true
process_name=mwopnet
