# s6/services.d/swoole/run
#!/usr/bin/execlineb -P
s6-setuidgid zendphp
/etc/services.d/swoole/daemon
