# s6/services.d/swoole/run
#!/usr/bin/with-contenv sh
s6-setuidgid zendphp /etc/services.d/swoole/daemon
