# .docker/s6/services.d/supervisor/run
#!/usr/bin/with-contenv sh
/etc/services.d/supervisor/daemon
