# s6/services.d/cron/run
#!/bin/bash
# From https://stackoverflow.com/a/51442053
cron && touch /var/log/cron.log && tail -F /var/log/cron.log
